# ops/cron/kyo-no-kyoto-crawl.cron.example
# Kyo no Kyoto daily crawl cycle
# Runs every day at 03:15 server time.
15 3 * * * cd /srv/kyo-no-kyoto && export PATH="$HOME/.nvm/versions/node/v22.22.0/bin:$PATH" && /usr/bin/env node scripts/run-crawl-cycle.mjs >> /var/log/kyo-no-kyoto-crawl.log 2>&1
