# ops/systemd/kyo-no-kyoto-crawl@.service.example
[Unit]
Description=Run Kyo no Kyoto crawl cycle for %i
Wants=network-online.target
After=network-online.target
OnFailure=kyo-no-kyoto-crawl-failure@%i.service

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/srv/kyo-no-kyoto
Environment=PATH=/home/ubuntu/.nvm/versions/node/v22.22.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=CRAWL4AI_PYTHON=/srv/kyo-no-kyoto/apps/crawler/.venv/bin/python
Environment=CRAWL4_AI_BASE_DIRECTORY=/srv/kyo-no-kyoto/apps/crawler/.cache
ExecStart=/usr/bin/flock -w 10800 /run/lock/kyo-no-kyoto-crawl.lock /usr/bin/env node scripts/run-crawl-cycle.mjs --city=%i --trigger=scheduled
StandardOutput=append:/var/log/kyo-no-kyoto-crawl-%i.log
StandardError=append:/var/log/kyo-no-kyoto-crawl-%i.log
