# ops/systemd/kyo-no-kyoto-crawl@tokyo.timer.example
[Unit]
Description=Run Tokyo crawl cycle every 36 hours

[Timer]
OnActiveSec=5h10m
OnUnitActiveSec=36h
Unit=kyo-no-kyoto-crawl@tokyo.service

[Install]
WantedBy=timers.target
