# ops/systemd/kyo-no-kyoto-crawl-failure@.service.example
[Unit]
Description=Report failed Kyo no Kyoto crawl cycle for %i

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/srv/kyo-no-kyoto
Environment=PATH=/home/ubuntu/.nvm/versions/node/v22.22.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env node scripts/report-crawl-failure.mjs --city=%i
