# ops/systemd/kyo-no-kyoto-crawl@tokyo.timer.example
[Unit]
Description=Run Tokyo crawl cycle every 36 hours

[Timer]
OnBootSec=5h10m
OnUnitActiveSec=36h
Persistent=true
Unit=kyo-no-kyoto-crawl@tokyo.service

[Install]
WantedBy=timers.target
