# ops/systemd/kyo-no-kyoto-crawl@kyoto.timer.example
[Unit]
Description=Run Kyoto crawl cycle every 36 hours

[Timer]
OnActiveSec=10m
OnUnitActiveSec=36h
Unit=kyo-no-kyoto-crawl@kyoto.service

[Install]
WantedBy=timers.target
